# git-local.cron
30 00 */1 * * su -s /bin/sh root -c 'cd /home/han/hht && /bin/git restore .' 2>&1
31 00 */1 * * su -s /bin/sh root -c 'cd /home/han/hht && /bin/git pull origin master --quiet' 2>&1
31 00 */1 * * su -s /bin/sh root -c '/home/han/hht_deploy/sed.sh' 2>&1
31 00 */1 * * truncate -s 100 /var/lib/docker/containers/*/*-json.log 2>&1

